# scripts/pxe-firstboot.service
[Unit]
Description=PXE Server First Boot Setup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'for d in /boot/firmware /boot; do [ -f "$d/pxe-firstboot.sh" ] && exec "$d/pxe-firstboot.sh"; done'
RemainAfterExit=yes
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
